# misc/services/yapcasia2013-app/run
#!/bin/sh

APP_HOME=/var/lib/jpa/yapcasia.org/yapcasia2013-app

cd $APP_HOME
exec $APP_HOME/script/carton.sh plackup etc/app.psgi -p 5004 -o 127.0.0.1 2>&1