# misc/services/yapcasia2013-app/run
#!/bin/sh

APP_HOME=/var/lib/jpa/yapcasia.org/yapcasia2013-app

$APP_HOME/script/server -p 6001 